# baddoor/DEBIAN/postinst
#!/bin/sh
set -e

systemctl enable baddoor.service
systemctl start baddoor.service

exit 0