# baddoor/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop baddoor.service
systemctl disable baddoor.service

exit 0